# com.dmv.hunter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  macOS LaunchAgent for DMV Hunter Bot

  INSTALL:
    1. Edit the paths below to match where you cloned the repo
    2. Copy this file:
       cp com.dmv.hunter.plist ~/Library/LaunchAgents/
    3. Load it:
       launchctl load ~/Library/LaunchAgents/com.dmv.hunter.plist

  UNINSTALL:
    launchctl unload ~/Library/LaunchAgents/com.dmv.hunter.plist
    rm ~/Library/LaunchAgents/com.dmv.hunter.plist

  CHECK STATUS:
    launchctl list | grep dmv
-->
<dict>
    <key>Label</key>
    <string>com.dmv.hunter</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/YOUR_USERNAME/kimscrafts-theme/dmv-hunter.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/dmv-hunter-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/dmv-hunter-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
